# resources/systemd_units/selenium_hub.service
[Unit]
Description=Selenium Hub
After=network.target

[Service]
ExecStart=/usr/lib/selenium/selenium_hub.sh

[Install]
WantedBy=multi-user.target